# roles/strategist/scripts/systemd/iwe-strategist-morning.service
[Unit]
Description=IWE Strategist — morning run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{IWE_RUNTIME}}/roles/strategist/scripts/strategist.sh morning
Environment=HOME={{HOME_DIR}}
Environment=IWE_TEMPLATE={{IWE_TEMPLATE}}
Environment=IWE_WORKSPACE={{WORKSPACE_DIR}}
Environment=IWE_RUNTIME={{IWE_RUNTIME}}
Environment=IWE_GOVERNANCE_REPO={{GOVERNANCE_REPO}}
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
StandardOutput=append:{{HOME_DIR}}/logs/strategist/systemd-morning.log
StandardError=append:{{HOME_DIR}}/logs/strategist/systemd-morning-error.log
